# templates/default/cassandra.service.erb
[Unit]
Description=cassandra daemon
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=<%= node['cassandra']['startup_program'] %> -f -p <%= node['cassandra']['pid_dir'] %>/cassandra.pid
Type=simple
User=<%= node['cassandra']['user'] %>
Group=<%= node['cassandra']['group'] %>
StandardOutput=journal
StandardError=inherit
Environment="CASSANDRA_HOME=<%= node['cassandra']['installation_dir'] %>" "CASSANDRA_CONF=<%= node['cassandra']['conf_dir'] %>"
LimitNOFILE=infinity
TimeoutStopSec=10min
Restart=<%= node['cassandra']['systemd']['restart_condition'] %>
PIDFile=<%= node['cassandra']['pid_dir'] %>/cassandra.pid
ExecStop=<%= node['cassandra']['bin_dir'] %>/nodetool  disablegossip
ExecStop=<%= node['cassandra']['bin_dir'] %>/nodetool  disablethrift
ExecStop=<%= node['cassandra']['bin_dir'] %>/nodetool  drain
ExecStop=-<%= node['cassandra']['bin_dir'] %>/nodetool stopdaemon
